# scripts/boot/99-coralreef-permissions.rules
# coralReef udev rules — passwordless GPU access for biomegate
#
# Sets permissions so all VFIO and GPU sysfs operations work
# without sudo/pkexec. Install to /etc/udev/rules.d/

# VFIO devices: world read/write (container + group fds)
SUBSYSTEM=="vfio", MODE="0666"

# PCI sysfs attributes for coralReef-managed Titan V cards:
# Allow biomegate to write driver_override, power/control, config
# Target: 0000:4a:00.0 (VFIO target) and 0000:03:00.0 (oracle)

# GPU function 0 - driver override (for hot-swap)
ACTION=="add|change", SUBSYSTEM=="pci", KERNEL=="0000:4a:00.0", \
  ATTR{vendor}=="0x10de", ATTR{device}=="0x1d81", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:4a:00.0/driver_override", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:4a:00.0/config", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:4a:00.0/resource0"

ACTION=="add|change", SUBSYSTEM=="pci", KERNEL=="0000:03:00.0", \
  ATTR{vendor}=="0x10de", ATTR{device}=="0x1d81", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:03:00.0/driver_override", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:03:00.0/config", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:03:00.0/resource0"

# Audio function - driver override
ACTION=="add|change", SUBSYSTEM=="pci", KERNEL=="0000:4a:00.1", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:4a:00.1/driver_override"

ACTION=="add|change", SUBSYSTEM=="pci", KERNEL=="0000:03:00.1", \
  RUN+="/bin/chmod 0666 /sys/bus/pci/devices/0000:03:00.1/driver_override"
