# scripts/boot/99-coralreef-vfio.rules
# coralReef udev rules for per-device driver binding at boot
#
# Fleet: Titan V #1 (03:00.0), RTX 5070 (21:00.0), Tesla K80 (4c:00.0 + 4d:00.0)
# Titan V #1: vfio-pci (oracle, swappable to nouveau via ember)
# RTX 5070: nvidia (display, locked — no rule needed, nvidia claims by default)
# K80 die #1 (4c:00.0): vfio-pci (sovereign target, cold at boot)
# K80 die #2 (4d:00.0): vfio-pci (reagent target, cold at boot)
#
# K80 DESIGN NOTE: Both K80 dies boot on vfio-pci (cold). Nouveau cannot
# cold-start GK210 — BOOT0 returns 0x0f22d0a1 (hardware default, not a valid
# chipset ID) because UEFI only POSTs the primary display GPU.
#
# Initialization path: ember detects cold hardware via BOOT0/PTIMER heuristics,
# then the K80 is POSTed by passing it to an agentReagents micro-VM running
# nvidia-470. After POST, ember swaps to vfio-pci for sovereign compute.
#
# IMPORTANT: K80 (10de:102d) should NOT be in the kernel cmdline vfio-pci.ids.
# Ember manages K80 binding via driver_override + drivers_autoprobe control.
# Only Titan V (10de:1d81,10de:10f2) needs cmdline protection because nvidia's
# open module corrupts Volta HBM2 state on failed probe.
#
# Install:
#   sudo cp scripts/boot/99-coralreef-vfio.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules

# --- Tesla K80 die #1 -> vfio-pci (sovereign target, cold) ---
ACTION=="add", KERNEL=="0000:4c:00.0", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", ATTR{device}=="0x102d", RUN+="/bin/sh -c 'echo vfio-pci > /sys/bus/pci/devices/0000:4c:00.0/driver_override'"
ACTION=="add", KERNEL=="0000:4c:00.0", SUBSYSTEM=="pci", ATTR{power/control}="on"
ACTION=="add", KERNEL=="0000:4c:00.0", SUBSYSTEM=="pci", RUN+="/bin/sh -c 'echo > /sys/bus/pci/devices/0000:4c:00.0/reset_method 2>/dev/null || true'"
ACTION=="bind", KERNEL=="0000:4c:00.0", SUBSYSTEM=="pci", RUN+="/bin/sh -c 'sleep 1 && chmod 666 /dev/vfio/37 2>/dev/null || true'"

# --- Tesla K80 die #2 -> vfio-pci (reagent target, cold) ---
ACTION=="add", KERNEL=="0000:4d:00.0", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", ATTR{device}=="0x102d", RUN+="/bin/sh -c 'echo vfio-pci > /sys/bus/pci/devices/0000:4d:00.0/driver_override'"
ACTION=="add", KERNEL=="0000:4d:00.0", SUBSYSTEM=="pci", ATTR{power/control}="on"
ACTION=="add", KERNEL=="0000:4d:00.0", SUBSYSTEM=="pci", RUN+="/bin/sh -c 'echo > /sys/bus/pci/devices/0000:4d:00.0/reset_method 2>/dev/null || true'"
ACTION=="bind", KERNEL=="0000:4d:00.0", SUBSYSTEM=="pci", RUN+="/bin/sh -c 'sleep 1 && chmod 666 /dev/vfio/38 2>/dev/null || true'"

# --- Titan V #1 -> vfio-pci (oracle, swappable to nouveau via ember) ---
ACTION=="add", KERNEL=="0000:03:00.0", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", ATTR{device}=="0x1d81", RUN+="/bin/sh -c 'echo vfio-pci > /sys/bus/pci/devices/0000:03:00.0/driver_override'"
ACTION=="add", KERNEL=="0000:03:00.0", SUBSYSTEM=="pci", ATTR{power/control}="on"
ACTION=="bind", KERNEL=="0000:03:00.0", SUBSYSTEM=="pci", RUN+="/bin/sh -c 'sleep 1 && chmod 666 /dev/vfio/73 2>/dev/null || true'"

# Titan V #1 audio -> vfio-pci (same IOMMU group)
ACTION=="add", KERNEL=="0000:03:00.1", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", RUN+="/bin/sh -c 'echo vfio-pci > /sys/bus/pci/devices/0000:03:00.1/driver_override'"

# VFIO devices: world read/write
SUBSYSTEM=="vfio", MODE="0666"
