# scripts/boot/coral-glowplug.service
[Unit]
Description=coral-glowplug — Sovereign PCIe Device Lifecycle Broker
Documentation=https://github.com/ecoPrimals/coralReef
After=local-fs.target systemd-modules-load.service coral-ember.service
Wants=systemd-modules-load.service coral-ember.service
Before=display-manager.service

[Service]
Type=notify
NotifyAccess=main
ExecStartPre=/bin/sh -c 'echo Y > /sys/module/vfio_pci/parameters/disable_idle_d3 2>/dev/null || true'
ExecStart=/usr/local/bin/coral-glowplug --config /etc/coralreef/glowplug.toml
Restart=on-failure
RestartSec=5
TimeoutStartSec=90
TimeoutStopSec=8
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
WatchdogSec=120
Environment=RUST_LOG=info
# RuntimeDirectory owned by coral-ember.service (ember starts first, never stops).
# Glowplug must NOT declare RuntimeDirectory — systemd would wipe /run/coralreef/
# on glowplug restart, destroying ember.sock and breaking SCM_RIGHTS fd passing.

[Install]
WantedBy=multi-user.target
