# scripts/boot/coral-glowplug.service
[Unit]
Description=coralReef GlowPlug — Sovereign PCIe Device Lifecycle Broker
After=coral-ember.service
Wants=coral-ember.service

[Service]
Type=notify
NotifyAccess=all
ExecStart=/usr/local/bin/coral-glowplug --config /etc/coralreef/glowplug.toml
RuntimeDirectory=coralreef
RuntimeDirectoryMode=0755
Environment=RUST_LOG=info

CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_DAC_OVERRIDE
AmbientCapabilities=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_DAC_OVERRIDE
NoNewPrivileges=true

ProtectSystem=strict
ReadWritePaths=/sys/bus/pci /run/coralreef /proc
ReadOnlyPaths=/etc/coralreef /sys/kernel/iommu_groups
PrivateTmp=true
ProtectHome=true

SystemCallFilter=@system-service ioctl recvmsg
SystemCallErrorNumber=EPERM
MemoryDenyWriteExecute=true

Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
