# scripts/boot/coral-ember.service
[Unit]
Description=coral-ember — Immortal VFIO fd Holder
Documentation=https://github.com/ecoPrimals/coralReef
After=local-fs.target systemd-modules-load.service
Wants=systemd-modules-load.service
Before=coral-glowplug.service display-manager.service

[Service]
Type=notify
NotifyAccess=main
ExecStartPre=/bin/sh -c 'echo Y > /sys/module/vfio_pci/parameters/disable_idle_d3 2>/dev/null || true'
ExecStart=/usr/local/bin/coral-ember /etc/coralreef/glowplug.toml
# CRITICAL: Restart=no — if ember dies, VFIO fds are lost and PM reset
# will occur. Only a cold reboot can recover. The ember must never crash.
Restart=no
TimeoutStartSec=30
TimeoutStopSec=3
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
Environment=RUST_LOG=info
Environment=CORALREEF_EMBER_SOCKET=/run/coralreef/ember.sock
SupplementaryGroups=coralreef
UMask=0117
RuntimeDirectory=coralreef
RuntimeDirectoryMode=0750
ExecStartPost=/bin/chgrp coralreef /run/coralreef

CapabilityBoundingSet=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_DAC_OVERRIDE CAP_SYS_MODULE
AmbientCapabilities=CAP_SYS_ADMIN CAP_SYS_RAWIO CAP_DAC_OVERRIDE CAP_SYS_MODULE
NoNewPrivileges=true

ProtectSystem=strict
ReadWritePaths=/sys/bus/pci /sys/kernel/iommu_groups /run/coralreef /proc /sys/kernel/debug/tracing /var/lib/coralreef
ReadOnlyPaths=/lib/modules /etc/modprobe.d /lib/modprobe.d
PrivateTmp=true
ProtectHome=true

SystemCallFilter=@system-service @module ioctl sendmsg
SystemCallErrorNumber=EPERM
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
